# systemd/user/ac.target
[Unit]
Description=User level ac.target
DefaultDependencies=no
StopWhenUnneeded=yes
