# systemd/user/battery.target
[Unit]
Description=User level batter.target
DefaultDependencies=no
StopWhenUnneeded=yes
