# systemd/user/widget-update.service
[Unit]
Description=Update qtile widget

[Service]
Type=oneshot
ExecStart=/usr/local/bin/widget-update

[Install]
WantedBy=battery.target ac.target
